# deploy/systemd/beichen-alpha-data-prewarm.service
[Unit]
Description=Beichen Alpha 30-minute data prewarm
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/www/wwwroot/daozang-alpha-suite/beichen-alpha
Environment=PYTHON_BIN=/www/wwwroot/daozang-alpha-suite/beichen-alpha/.venv/bin/python
Environment=OMP_NUM_THREADS=1
Environment=OPENBLAS_NUM_THREADS=1
Environment=MKL_NUM_THREADS=1
Environment=NUMEXPR_NUM_THREADS=1
ExecStart=/www/wwwroot/daozang-alpha-suite/beichen-alpha/scripts/beichen_data_prewarm.sh
Nice=10
CPUQuota=40%
MemoryMax=900M
TasksMax=64
IOWeight=50
