# deploy/systemd/beichen-alpha-data-prewarm.timer
[Unit]
Description=Run Beichen Alpha data prewarm every 30 minutes during A-share sessions

[Timer]
OnCalendar=Mon..Fri *-*-* 09:30:00
OnCalendar=Mon..Fri *-*-* 10:00:00
OnCalendar=Mon..Fri *-*-* 10:30:00
OnCalendar=Mon..Fri *-*-* 11:00:00
OnCalendar=Mon..Fri *-*-* 11:30:00
OnCalendar=Mon..Fri *-*-* 13:00:00
OnCalendar=Mon..Fri *-*-* 13:30:00
OnCalendar=Mon..Fri *-*-* 14:00:00
OnCalendar=Mon..Fri *-*-* 14:30:00
OnCalendar=Mon..Fri *-*-* 15:00:00
RandomizedDelaySec=30
Persistent=false
Unit=beichen-alpha-data-prewarm.service

[Install]
WantedBy=timers.target
