# deploy/systemd/beichen-alpha-chat.service
[Unit]
Description=Beichen Alpha Feishu chat adapter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/www/wwwroot/daozang-alpha-suite/beichen-alpha
Environment=PYTHON_BIN=/www/wwwroot/daozang-alpha-suite/beichen-alpha/.venv/bin/python
Environment=PYTHONUNBUFFERED=1
ExecStart=/www/wwwroot/daozang-alpha-suite/beichen-alpha/scripts/beichen_chat_server.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
