# deploy/systemd/beichen-alpha.service
[Unit]
Description=Beichen Alpha daily trading assistant run
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/www/wwwroot/daozang-alpha-suite/beichen-alpha
Environment=PYTHON_BIN=/www/wwwroot/daozang-alpha-suite/beichen-alpha/.venv/bin/python
ExecStart=/www/wwwroot/daozang-alpha-suite/beichen-alpha/scripts/server_daily_run.sh
Nice=5
